# templates/check-update.cron
MAILTO=""
18 * * * * root /usr/local/bin/check-update.rb {{ update_check_snitch }}
